# stats-dashboard.service.example
[Unit]
Description=System Dashboard Flask App
After=network.target

[Service]
Type=simple
User=your-user
Group=your-user
WorkingDirectory=/opt/stats
Environment=PATH=/opt/stats/.venv/bin
ExecStart=/opt/stats/.venv/bin/gunicorn --workers 2 --bind 0.0.0.0:5000 app:app
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true

[Install]
WantedBy=multi-user.target
